# ATM_Jackpotting_DFIR_2025.json
{
  "name": "ATM Jackpotting DFIR 2025",
  "version": "4.5",
  "domain": "enterprise-attack",
  "description": "MITRE ATT&CK layer for ATM Jackpotting incidents, designed for DFIR in 2025. Covers techniques used by groups like Carbanak, Lazarus Group, Tren de Aragua, and unlisted groups, including malware (Ploutus, WinPot, Jackpot) and hardware (Black Box, SIM Card cable). Based on incidents in 2023-2024 (e.g., Thailand, US).",
  "techniques": [
    {
      "techniqueID": "T1078",
      "tactic": "initial-access",
      "color": "#ff0000",
      "comment": "Use of stolen master keys or endoscopes for physical access to ATMs (e.g., US 2024, Thailand 2023)."
    },
    {
      "techniqueID": "T1190",
      "tactic": "initial-access",
      "color": "#ff9900",
      "comment": "Exploiting vulnerabilities in ATM software with network access (e.g., outdated Windows XP Embedded)."
    },
    {
      "techniqueID": "T1059",
      "tactic": "execution",
      "color": "#ff0000",
      "comment": "Malware (e.g., Ploutus.D, WinPot, Jackpot) executes scripts to control cash dispenser via XFS API."
    },
    {
      "techniqueID": "T1053",
      "tactic": "execution",
      "color": "#ff9900",
      "comment": "Scheduled tasks to run malicious scripts for cash dispensing."
    },
    {
      "techniqueID": "T1547",
      "tactic": "persistence",
      "color": "#ff9900",
      "comment": "Malware sets autostart execution to maintain access on ATM reboot."
    },
    {
      "techniqueID": "T1548",
      "tactic": "privilege-escalation",
      "color": "#ff9900",
      "comment": "Switching ATM to supervisor mode to gain control."
    },
    {
      "techniqueID": "T1562",
      "tactic": "defense-evasion",
      "color": "#ff9900",
      "comment": "Disabling security alerts or USB port locks to evade detection."
    },
    {
      "techniqueID": "T1070",
      "tactic": "defense-evasion",
      "color": "#ff9900",
      "comment": "Deleting logs to hide attack traces (e.g., Ploutus malware)."
    },
    {
      "techniqueID": "T1071",
      "tactic": "command-and-control",
      "color": "#ff0000",
      "comment": "Malware communicates with C2 server via HTTP/HTTPS (e.g., Carbanak, Lazarus Group)."
    },
    {
      "techniqueID": "T1090",
      "tactic": "command-and-control",
      "color": "#ff0000",
      "comment": "Use of SIM Card cable or Black Box for remote control (e.g., Thailand 2023, US 2024)."
    },
    {
      "techniqueID": "T1219",
      "tactic": "command-and-control",
      "color": "#ff0000",
      "comment": "Remote control via SMS/Bluetooth (e.g., Ploutus.D in US 2024)."
    },
    {
      "techniqueID": "T1485",
      "tactic": "impact",
      "color": "#ff9900",
      "comment": "Deleting logs or data to hide attack traces."
    },
    {
      "techniqueID": "T1496",
      "tactic": "impact",
      "color": "#ff0000",
      "comment": "Hijacking ATM to dispense cash (Jackpotting)."
    }
  ],
  "gradient": {
    "colors": ["#ff6666", "#ff0000"],
    "minValue": 0,
    "maxValue": 100
  },
  "legendItems": [
    {
      "color": "#ff0000",
      "label": "High-frequency technique"
    },
    {
      "color": "#ff9900",
      "label": "Supporting technique"
    }
  ]
}